# config/usbip-check.timer
[Unit]
Description=Periodically Check USB/IP Status
After=network.target

[Timer]
OnBootSec=60
OnUnitActiveSec=30s

[Install]
WantedBy=timers.target
